# srcs/nginx/supervisord.conf
[supervisord]
user = root
nodaemon=true

[program:sshd]
user = root
command = sh -c "/usr/sbin/sshd -D"
autostart=true
autorestart = true
startretries=10

[program:nginx]
user = root
command = sh -c "nginx -g \"daemon off;\""
autostart=true
autorestart = true
startretries=10

[program:telegraf]
user = root
command = sh -c "telegraf"
autostart=true
autorestart = true
startretries=10

