# srcs/pma/supervisord.conf
[supervisord]
user = root
nodaemon=true

[program:telegraf]
user = root
command = sh -c "telegraf"
autostart=true
autorestart = true
startretries=10

[program:nginx]
user = root
command = sh -c "nginx -g \"daemon off;\""
autostart=true
autorestart = true
startretries=10

[program:fpm]
user = root
command = sh -c "php-fpm7 --nodaemonize"
autostart=true
autorestart = true
startretries=10


